# configuration/linux/supervisor/vpnc.ini
[program:vpnc]
environment=HOME="/root",USER="root"
command=/usr/sbin/vpnc --local-port 10000 --enable-1des --debug 1 --no-detach
directory=/root/
user=root
autostart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/vpnc.log
